# vendors/Arcturus/uCquicc/inittab
inet:unknown:/bin/inetd
boa:unknown:/bin/boa
ttyS0:vt100:/bin/agetty -l /bin/login ttyS0 19200
