# vendors/Cwlinux/cwvz328/inittab
# inittab for uClinux
# Format:
# ttyline:termcap-entry:getty-command
inetd:unknown:/bin/inetd
boa:unknown:/bin/boa
ttyS0:vt100:/bin/agetty 115200 /dev/ttyS0
#ttyS0:vt100:/bin/sh
#ttyS1:vt100:/bin/sh
