# vendors/Netburner/MOD5272/inittab
# inittab for uClinux
# Format:
# ttyline:termcap-entry:getty-command
# ttyS0:vt100:/bin/agetty 9600 ttyS0
ttyS1:vt100:/bin/agetty 9600 ttyS1
inet:unknown:/bin/inetd
