# user/boa/debian/cron.daily
#!/bin/sh

savelog -c 7 /var/log/boa/access_log >/dev/null
savelog -c 7 /var/log/boa/error_log >/dev/null 
/etc/init.d/boa force-reload > /dev/null
