# vendors/Simtec/EB675001DIP/inittab
inet:unknown:/bin/inetd
boa:unknown:/bin/boa -c /etc/boa
