# vendors/McAfee/audiophile/romfs/etc/inittab
flt:unknown:/bin/flatfsd
mb:unknown:/bin/mb
cron:unknown:/bin/cron
agetty:unknown:/bin/agetty /dev/ttyS0 115200
