# vendors/Insight/romfs/etc/default/inittab
boa:unknown:/bin/boa
dns:unknown:/bin/dnsmasq -d
