# vendors/Simtec/EB675001DIP-ssh/inittab
inet:unknown:/bin/inetd
boa:unknown:/bin/boa -c /etc/boa
ttyS0:vt100:/bin/agetty 115200 ttyS0
