# vendors/Soekris/net4801/inittab
::sysinit:/bin/sh /etc/rc
::respawn:/sbin/syslogd -n
::respawn:/sbin/klogd -n
::respawn:/bin/inetd
::respawn:/bin/flatfsd
