# vendors/TeleIP/rsc/inittab
# inittab for uClinux
#inet:unknown:/bin/inetd
#boa:unknown:/bin/boa

# Format:
# ttyline:termcap-entry:getty-command
#ttyS0:vt100:/sbin/agetty 19200 ttyS0
#ttyS0:vt100:/bin/getty -iL 19200 ttyS0
ttyS0:vt100:/bin/sh
