# vendors/Weiss/SM2010/inittab
#ttyS1:linux:/bin/sh
#pppd listen on /dev/ttyS1
#s1:linux:/bin/pppd /dev/ttyS1 57600 noauth noipdefault defaultroute
