# vendors/EMAC/SOM5282EM/inittab
inet:unknown:/bin/inetd
boa:unknown:/bin/boa
/dev/console:vt100:/sbin/getty 115200 /dev/console vt100
