# vendors/Serra/big/inittab
inet:unknown:/bin/inetd
boa:unknown:/bin/boa
nmb:unknown:/bin/nmbd -s /home/conf/smb.conf
#smb:unknown:/bin/smbd -s /home/conf/smb.conf
