# vendors/Signal/romfs/etc/inittab
console:unknown:/bin/sh /etc/rc_signal
inet:unknown:/bin/inetd
boa:unknown:/bin/boa
