# vendors/Arcturus/uC5275-16E32U150/inittab
# inittab for uClinux
# Format:
# ttyline:termcap-entry:getty-command
#  To use the following line with speed 0 you need to use
#  agetty from the uCagetty package
#  uCagetty's implemenation retrieves the speed from the
#  bootloader environment variable CONSOLE_SPEED
#ttyS0:vt100:/bin/agetty 0 ttyS0
ttyS0:vt100:/bin/agetty 9600 ttyS0
ttyS1:vt100:/bin/agetty 9600 ttyS1
inet:unknown:/bin/inetd
http:unknown:/bin/thttpd -p 80 -d /home/httpd -c **.cgi -u root -e 5000
