# vendors/Xilinx/uclinux-auto/inittab
inet:unknown:/bin/inetd
ttyS1:unknown:/bin/getty ttyS1 115200
boa:unknown:/bin/boa

