# vendors/Insight/romfs/etc/default-static/inittab
boa:unknown:/bin/boa
dns:unknown:/bin/dnsmasq -d -l /etc/config/dhcpd.leases
dhsd:unknown:/bin/dhcpd eth0
