# vendors/Insight/XC2VP7_Microblaze/inittab
#inet:unknown:/bin/inetd
#flt:unknown:/bin/flatfsd
ttyS1:unknown:/bin/agetty ttyS1 57600

