# vendors/OpenGear/romfs/etc/inittab
inet:unknown:/bin/inetd
fltd:unknown:/bin/flatfsd
