# vendors/Insight/XC2V1K_Microblaze/inittab
#inet:unknown:/bin/inetd
#flt:unknown:/bin/flatfsd
ttyS1:unknown:/bin/getty ttyS1 115200

