# vendors/Philips/LPC2200/inittab
# inet:unknown:/bin/inetd
# boa:unknown:/bin/boa
