# vendors/Feith/CANCam/inittab
inet:unknown:/bin/inetd
flt:unknown:/bin/flatfsd
