# vendors/SecureEdge/eLIA/inittab
inet:unknown:/bin/inetd
boa:unknown:/bin/boa
flt:unknown:/bin/flatfsd
