# vendors/Apple/iPod/inittab
inet:unknown:/bin/inetd
pz:unknown:/bin/podzilla

