# vendors/SecureEdge/SecureEdgeEdna/romfs/etc/inittab
inet:unknown:/bin/inetd
flt:unknown:/bin/flatfsd
pmap:unknown:/bin/portmap
